# etc/cron.daily/borgmatic
#!/bin/bash

su tomeraberbach
. ~/.keychain/$HOSTNAME-sh

export BORG_PASSCOMMAND='sudo -u tomeraberbach pass borgbase.com/passphrase'

borgmatic --syslog-verbosity 2

