# home/tomeraberbach/.config/autostart/unlock-secrets.desktop
[Desktop Entry]
Exec=/usr/bin/unlock-secrets
Name=Unlock Secrets
Icon=system-run
StartupNotify=true
Terminal=true
Type=Application
